# uudev.service.in
[Unit]
Description=Micro udev service

[Service]
Type=simple
ExecStart=@BINDIR@/bin/uudev

[Install]
WantedBy=default.target
